# padelizou.service
[Unit]
Description=Padelizou - App de torneios e aulas de padel
# O banco fica no Azure SQL (remoto) - só precisa da rede local, não de um SQL Server na mesma máquina.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/padelizou
ExecStart=/usr/bin/dotnet /opt/padelizou/Padelizou.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=padelizou
User=www-data

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

Environment="ConnectionStrings__DefaultConnection=Server=tcp:SEUSERVIDOR.database.windows.net,1433;Database=DB_PADEL;User Id=SEU_LOGIN_ADMIN;Password=SENHA_FORTE_AQUI;Encrypt=True;TrustServerCertificate=False;Connection Timeout=30;"
Environment="EmailSettings__SmtpHost=smtp.gmail.com"
Environment="EmailSettings__SmtpPort=587"
Environment="EmailSettings__RemetenteEmail=seuemail@gmail.com"
Environment="EmailSettings__RemetenteSenhaApp=senha-de-app-do-gmail"
Environment="EmailSettings__RemetenteNome=Padelizou"
Environment="GoogleCalendar__ClientId=SEU_CLIENT_ID_AQUI"
Environment="GoogleCalendar__ClientSecret=SEU_CLIENT_SECRET_AQUI"
Environment="GoogleCalendar__RedirectUri=https://padelizou.com.br/GoogleAuth/Callback"
Environment="AcessoAntecipado__Habilitado=true"
Environment="AcessoAntecipado__Usuario=padelizou"
Environment="AcessoAntecipado__Senha=TROCAR_ANTES_DE_PUBLICAR"

[Install]
WantedBy=multi-user.target
